# circuit-bot.timer.example
[Unit]
Description=Periodically run Circuit bot for status notifications

[Timer]
OnCalendar=minutely
